# .config/supervisor/uwsgi.conf
[program:uwsgi]
command=/root/.pyenv/versions/weathersound/bin/uwsgi --ini /srv/weathersound/.config/uwsgi/uwsgi-app.ini

